# ui/install/ui.conf
[program:ui]
command=python3 app.py
directory=/home/pi/streamer-pi/ui
autostart=true
autorestart=true
startretries=100000
stdout_logfile=/home/pi/log/ui.log
stderr_logfile=/home/pi/log/ui_error.log

